# UpdateOS-3.1/system/lib/systemd/system/createap.service
[Unit] 
Description=Create AP Daemon

[Service] 
Type=simple 
WorkingDirectory=/root/
ExecStart=/root/scripts/start-ap.sh

[Install] 
WantedBy=multi-user.target
